# deploy/systemd/kis-kill-switch.service
[Unit]
Description=KIS Kill Switch (Phase 4, 6-condition monitor)
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/project/kis_unified_sts
EnvironmentFile=/home/deploy/project/kis_unified_sts/.env
ExecStart=/home/deploy/project/kis_unified_sts/.venv/bin/python -m services.kill_switch.main
# When the kill switch trips it exits the run loop intentionally — do NOT
# auto-restart; an operator must run scripts/kill_switch_clear.sh first.
Restart=no
TimeoutStopSec=30s
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
