# deploy/systemd/kis-stock-monitor-daemon.service
[Unit]
Description=KIS Stock Monitor Bridge (daemon streams -> dashboard state + alerts)
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/project/kis_unified_sts
EnvironmentFile=/home/deploy/project/kis_unified_sts/.env
Environment=STOCK_MONITOR_DAEMON=shadow
Environment=STOCK_POSITIONS_KEY=stock:daemon:positions
Environment=TRADING_STATE_KEY_SUFFIX=shadow
ExecStart=/home/deploy/project/kis_unified_sts/.venv/bin/python -m services.stock_monitor.main
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s
KillSignal=SIGTERM

# Delivered DISABLED. Enabling is an operator step (shadow validation gate).
[Install]
WantedBy=multi-user.target
