# deploy/scheduler.crontab
# KIS scheduler crontab — run by the `scheduler` compose service via supercronic.
# Times are KST (the service sets TZ=Asia/Seoul; supercronic honours TZ).
# Replaces the host crontab KIS jobs (cron→compose migration,
# docs/plans/archive/2026-06-09-cron-to-compose-scheduler.md). One-shot jobs only — the
# market-hours daemons (screener/fusion) are separate compose services.
#
# Entrypoints run from /app (the app image); no host .sh wrappers, no `sts`
# console-script (PYTHONPATH=/app is set on the service for `-m cli.main`).

# --- pre-market (06:30) ---
30 6  * * 1-5  cd /app && python scripts/llm_premarket_briefing.py
30 6  * * 1-5  cd /app && python -m services.macro_overnight_collector.main us

# --- intraday macro (every 15 min) ---
*/15 * * * 1-5  cd /app && python -m services.macro_overnight_collector.main fx

# --- pre-market macro (07:45) ---
# ES=F/NQ=F/^SOX + offshore KRW=X for the 08:00 pre-open Risk Score.
# ECOS publishes the same-day USD/KRW fixing only after 08:30 KST, so the
# realtime offshore rate bridges the gap. Symbols come from
# config/macro_sources.yaml (yahoo_symbols); keep its sessions doc-block in
# sync with this schedule.
45 7  * * 1-5  cd /app && python -m services.macro_overnight_collector.main premarket

# --- forecasting refit (07:35) ---
# --from-parquet: auto-resolves active A01* near-month from Parquet minute bars,
# fits HAR-RV, and writes forecast:vol:model to Redis.  No external DB required.
35 7  * * 1-5  cd /app && python scripts/forecasting/refit_har_rv.py --from-parquet

# --- pre-open scanners (08:30 / 08:50 / 08:58) → daily_watchlist, indicators ---
30 8  * * 1-5  cd /app && python scripts/run_daily_scanner.py
50 8  * * 1-5  cd /app && python scripts/daily_indicator_scanner.py
58 8  * * 1-5  cd /app && python scripts/daily_indicator_scanner.py

# --- pre-open risk reset (08:59) ---
# Idempotently resets decoupled stock/futures daily risk counters before
# risk-filter daemons admit the first regular-session entries.
59 8  * * 1-5  cd /app && python -m scripts.maintenance.daily_risk_reset

# --- intraday LLM refresh (hourly 10:00–15:00) ---
0 10-15 * * 1-5  cd /app && python -m scripts.analysis.llm_intraday_refresh

# --- market close (15:30 / 15:32 / 15:40) ---
30 15 * * 1-5  cd /app && python -m scripts.analysis.llm_market_close_briefing
# Setup A/C paper-observation daily digest — runs 2 min after the main briefing
# so it appears as a follow-on message in the briefing channel.
32 15 * * 1-5  cd /app && python -m scripts.analysis.setup_ac_paper_observation >> /app/logs/setup_ac_paper_obs_$(date +\%Y\%m\%d).log 2>&1
40 15 * * 1-5  cd /app && python scripts/analysis/publish_equity_snapshot.py
# LLM scorecard post-close scorer — 16:35 KST after close+settlement
35 16 * * 1-5  cd /app && python -m scripts.analysis.llm_scorecard_score >> /app/logs/scorecard_$(date +\%Y\%m\%d).log 2>&1

# --- EOD backfill + verification (15:50–16:20) ---
# --days 30 keeps the full KIS-max minute window filled (heals failed days +
# denser intraday history for experiments); state-DB resume skips done (code,date)
# so it stays cheap. KIS caps 1-min history at ~30 days — can't go deeper.
50 15 * * 1-5  cd /app && python -m cli.main stock-backfill run --days 30
5  16 * * 1-5  cd /app && python -m cli.main backfill run
10 16 * * 1-5  cd /app && python scripts/analysis/stock_paper_daily_verification.py
# Deep daily backfill: --days 400 lets the paginating fetch fill >100 bars/code
# (KIS caps a single call at ~100), so SMA(200) symbols stay deep. State-DB resume
# keeps it cheap (only missing (code,date) are fetched).
20 16 * * 1-5  cd /app && python -m cli.main stock-backfill daily --days 400

# --- on-entry coverage worker (every 10 min, market hours 08:00–15:50 KST) ---
# Drains the Redis coverage queue (filled by market-ingest when a symbol enters
# the universe) and deepens any shallow symbol to >= 200 daily bars so newly
# admitted LLM/screener symbols are trade-ready for pattern_pullback the same
# session. Idempotent + throttled; no-op when the queue is empty.
*/10 8-15 * * 1-5  cd /app && python -m cli.main stock-backfill ensure-coverage

# Futures-minute coverage gate — Telegram alert if the active front-month is
# short vs expected bars (catches the ~70% loss class, PR #507). After 16:05 backfill.
30 16 * * 1-5  cd /app && python scripts/analysis/futures_minute_coverage_gate.py

# =============================================================================
# Continuous off-hours backfill — BOTH assets (data-engineer follow-up (b)).
# Backfill runs CONTINUOUSLY across
# the advantageous off-hours windows — 장전 (~07:00–08:50), 장후 (~15:40+), 야간
# (overnight) — and NEVER inside the live session 09:00–15:30 KST (so it never
# competes with the live feed). Each job is resume-based + gate-driven: state-DB
# resume skips done (code,date) and the completeness gate refetches only what is
# missing, so repeated passes are cheap and idempotent. KST-native (TZ=Asia/Seoul).
# KIS maintenance window (~23:40–00:10 KST) is avoided: night jobs cap at 23:00
# and the pre-open pass resumes at 07:10. All windows tagged below.
#
# IMAGE REBUILD REQUIRED: this crontab is baked into the app image; rebuild the
# scheduler service image and recreate the container for changes to take effect.
#
# Futures-full backfill (--futures): collects the A01* underlying (kospi200f_1m)
# AND rebuilds the 101S6000 continuous (kospi200f_continuous_1m). Pre-#518 only
# `mini` ran; this continuous schedule re-attempts the 101S6000 June=0 gap.
# --days 30 fills the full KIS-max minute window; resume keeps it cheap.
#
# --- 장후 (post-close, 15:40+) ---
35 16 * * 1-5  cd /app && python -m cli.main backfill run --days 30 --futures
# --- 야간 (overnight, hourly 19:00–23:00; capped before KIS maintenance) ---
0 19-23 * * 1-5  cd /app && python -m cli.main backfill run --days 30 --futures
# --- 장전 (pre-open, after 07:00 wake) ---
10 7  * * 1-5  cd /app && python -m cli.main backfill run --days 30 --futures

# Stock minute top-up — keep the KIS-max 30-day minute window dense for the live
# universe (heals failed days, denser intraday history). Resume skips done days.
# Staggered off :00 so the night pass never shares a minute with the hourly
# futures pass (0 19-23) and contends for the KIS rate limiter.
# --- 야간 ---
5 21 * * 1-5  cd /app && python -m cli.main stock-backfill run --days 30
# --- 장전 ---
20 7  * * 1-5  cd /app && python -m cli.main stock-backfill run --days 30

# Stock daily top-up — keep SMA(200)-deep daily history for the universe.
# Paginating fetch (>100 bars/code); resume keeps it cheap. Staggered off :00 so
# the night pass never collides with the hourly futures pass (0 19-23).
# --- 야간 ---
25 22 * * 1-5  cd /app && python -m cli.main stock-backfill daily --days 400

# --- accumulation pattern scan (16:10) → system:accumulation:latest ---
10 16 * * 1-5  cd /app && python scripts/run_accumulation_scan.py

# --- strategy experiment (16:40, after daily backfill) → reports/stock_experiment ---
# Backtests the current registry strategies on collected data via the real
# BacktestEngine; the dashboard /experiments page reads the latest report.
40 16 * * 1-5  cd /app && python -m cli.main experiment run --spec config/experiments/stock_default.yaml >> /app/logs/stock_experiment_$(date +\%Y\%m\%d).log 2>&1

# --- maintenance / weekly ---
0 4  * * 0  cd /app && python -m scripts.maintenance.rotate_reports
0 18 * * 0  cd /app && python scripts/analysis/regime_gate_counterfactual.py
# Setup A/C paper-observation WEEKLY rollup — Friday 16:30 KST after EOD tasks.
30 16 * * 5  cd /app && python -m scripts.analysis.setup_ac_paper_observation --weekly >> /app/logs/setup_ac_paper_obs_$(date +\%Y\%m\%d)_weekly.log 2>&1
# LLM scorecard weekly digest — Friday 16:37 KST (after daily scorer at 16:35).
37 16 * * 5  cd /app && python -m scripts.analysis.llm_scorecard_weekly >> /app/logs/scorecard_weekly_$(date +\%Y\%m\%d).log 2>&1

# --- DROPPED (PR-B): host crons that point at missing/stale scripts — NOT migrated ---
#   phase2_daily_verification  — scripts/cron/phase2_daily_verification.sh is MISSING (dead cron)
#   counterfactual_weekly      — scripts/cron/counterfactual_weekly.sh is MISSING (dead cron)
# PR-D host-cron teardown removes these dead entries from the host crontab.

# --- KRX night-derivatives close capture (05:48; window 05:50–06:00 KST, O9) ---
# The night session runs 18:00 → next-day 06:00 KST, so closes land Tue–Sat
# mornings (dow 2-6). One-shot WS (H0MFCNT0) capture of the last trade before
# the 06:00 close → Redis DB 1 market:structure:night_close (TTL 24h) for the
# pre-open Market Risk Score. Window/tr_key live in config/night_futures.yaml;
# starts 2 min early so connect/subscribe are warm before the window opens.
# Zero in-window trades (holiday/halt) → warning log, key not published.
48 5 * * 2-6  cd /app && python -m services.night_futures_collector.main

# --- market-structure daily snapshots (Wave 2b, unified roadmap Phase 0) ---
# docs/plans/2026-07-02-unified-investment-system-roadmap.md §Phase 0 + §4.3:
# premarket (08:00 KST) composes the pre-open knowledge row (prev close row +
# 07:45 macro premarket snapshot + night_close capture) — feeds the 08:00
# pre-open Risk Score without look-ahead; close (18:40 KST) captures the
# finalized daily row after the KRX ~18:00 investor/program publication
# (foreign futures FHPTJ04030000 snapshot ≥15:45 = daily confirmed value).
# Keys/TTLs/TR params live in config/market_structure.yaml (collector section).
0  8  * * 1-5  cd /app && python -m services.market_structure_collector.main premarket
40 18 * * 1-5  cd /app && python -m services.market_structure_collector.main close

# --- Market Risk Score engine (unified roadmap Phase 1, §4 — shadow-only) ---
# premarket (08:05, after the 08:00 structure collector) scores the pre-open
# row; intraday (*/30 09:00–15:30; session window filtered inside the runner
# via config/market_risk.yaml) re-scores market:structure:latest Redis-only;
# close (18:45, after the 18:40 close collector) is the confirmed daily
# computation — Redis + regime:unified:daily + score columns merged into the
# Parquet close row. Weights/bands/hysteresis/alerts: config/market_risk.yaml.
5  8  * * 1-5  cd /app && python -m services.market_risk_engine.main premarket
*/30 9-15 * * 1-5  cd /app && python -m services.market_risk_engine.main intraday
45 18 * * 1-5  cd /app && python -m services.market_risk_engine.main close

# --- Portfolio equity snapshot + unified monthly-MDD monitor (Phase 3B) ---
# 19:00 KST after the 18:45 market-risk close run. Track B/C equity =
# capital_base (config/portfolio.yaml) + RuntimeLedger track-tagged realized
# PnL + trading-state unrealized PnL; appends portfolio_equity_daily and
# publishes portfolio:equity:latest (TTL 24h) + stream:portfolio.equity.
# circuit_breaker.mode=shadow (default) → observe/alert only, no enforcement.
0 19 * * 1-5  cd /app && python -m services.portfolio_monitor.main
# Pre-market refresh (08:50 KST): re-publish portfolio:equity:latest before
# the session. Without it the Friday 19:00 publish's 24h TTL lapses Saturday
# 19:00, so Monday-intraday REDUCE/HALT_NEW enforcement would fail open all
# day. Idempotent: same one-shot, upsert by trade_date; the 19:00 run
# overwrites the row with the EOD state (risk-event audits dedupe by key).
50 8 * * 1-5  cd /app && python -m services.portfolio_monitor.main

# --- Unified performance feedback reports (Phase 6A — read-only analysis) ---
# Read-only batch: turns RuntimeLedger track-tagged trades/fills,
# portfolio_equity_daily, hedge_advice, market_structure close rows, and the
# newest experiment artifact into weekly/monthly/quarterly report files
# ({reports_root}/{weekly,monthly,quarterly}/<period_label>.{md,json}), a
# portfolio:feedback:latest Redis pointer, and one Telegram headline. It NEVER
# influences any strategy/execution/gate path (설계서 §8, §8.2). Defaults cover
# the just-completed period, so each run lands after that period closes.
# Weekly — Sat 09:00 KST: the trading week (Mon-Fri) has closed; the default
# period_label is that week's last trading day (Friday). After the Fri 19:00
# portfolio-monitor EOD row so the equity/slippage inputs are complete.
0 9 * * 6  cd /app && python -m services.feedback_reporter.main weekly
# Monthly — 1st 09:10 KST: covers the previous KST month (portfolio_equity_daily
# month rows + market-risk band residency are finalized by month end).
10 9 1 * *  cd /app && python -m services.feedback_reporter.main monthly
# Quarterly — Jan/Apr/Jul/Oct 1st 09:20 KST: §8.2 track judgment material for the
# just-closed quarter (produce-only; promotion/demotion stays with the operator).
20 9 1 1,4,7,10 *  cd /app && python -m services.feedback_reporter.main quarterly
