# deploy/scheduler.crontab
# KIS scheduler crontab — run by the `scheduler` compose service via supercronic.
# Times are KST (the service sets TZ=Asia/Seoul; supercronic honours TZ).
# Replaces the host crontab KIS jobs (cron→compose migration,
# docs/plans/archive/2026-06-09-cron-to-compose-scheduler.md). One-shot jobs only — the
# market-hours daemons (screener/fusion) are separate compose services.
#
# Entrypoints run from /app (the app image); no host .sh wrappers, no `sts`
# console-script (PYTHONPATH=/app is set on the service for `-m cli.main`).

# --- pre-market (06:30) ---
30 6  * * 1-5  cd /app && python scripts/llm_premarket_briefing.py
30 6  * * 1-5  cd /app && python -m services.macro_overnight_collector.main us

# --- intraday macro (every 15 min) ---
*/15 * * * 1-5  cd /app && python -m services.macro_overnight_collector.main fx

# --- forecasting refit (07:35) ---
35 7  * * 1-5  cd /app && python scripts/forecasting/refit_har_rv.py

# --- pre-open scanners (08:30 / 08:50 / 08:58) → daily_watchlist, indicators ---
30 8  * * 1-5  cd /app && python scripts/run_daily_scanner.py
50 8  * * 1-5  cd /app && python scripts/daily_indicator_scanner.py
58 8  * * 1-5  cd /app && python scripts/daily_indicator_scanner.py

# --- intraday LLM refresh (hourly 10:00–15:00) ---
0 10-15 * * 1-5  cd /app && python -m scripts.analysis.llm_intraday_refresh

# --- market close (15:30 / 15:32 / 15:40) ---
30 15 * * 1-5  cd /app && python -m scripts.analysis.llm_market_close_briefing
# Setup A/C paper-observation daily digest — runs 2 min after the main briefing
# so it appears as a follow-on message in the briefing channel.
32 15 * * 1-5  cd /app && python -m scripts.analysis.setup_ac_paper_observation >> /app/logs/setup_ac_paper_obs_$(date +\%Y\%m\%d).log 2>&1
40 15 * * 1-5  cd /app && python scripts/analysis/publish_equity_snapshot.py

# --- EOD backfill + verification (15:50–16:20) ---
# --days 30 keeps the full KIS-max minute window filled (heals failed days +
# denser intraday history for experiments); state-DB resume skips done (code,date)
# so it stays cheap. KIS caps 1-min history at ~30 days — can't go deeper.
50 15 * * 1-5  cd /app && python -m cli.main stock-backfill run --days 30
5  16 * * 1-5  cd /app && python -m cli.main backfill run
10 16 * * 1-5  cd /app && python scripts/analysis/stock_paper_daily_verification.py
20 16 * * 1-5  cd /app && python -m cli.main stock-backfill daily

# --- accumulation pattern scan (16:10) → system:accumulation:latest ---
10 16 * * 1-5  cd /app && python scripts/run_accumulation_scan.py

# --- strategy experiment (16:40, after daily backfill) → reports/stock_experiment ---
# Backtests the current registry strategies on collected data via the real
# BacktestEngine; the dashboard /experiments page reads the latest report.
40 16 * * 1-5  cd /app && python -m cli.main experiment run --spec config/experiments/stock_default.yaml >> /app/logs/stock_experiment_$(date +\%Y\%m\%d).log 2>&1

# --- maintenance / weekly ---
0 4  * * 0  cd /app && python -m scripts.maintenance.rotate_reports
0 18 * * 0  cd /app && python scripts/analysis/regime_gate_counterfactual.py
# Setup A/C paper-observation WEEKLY rollup — Friday 16:30 KST after EOD tasks.
30 16 * * 5  cd /app && python -m scripts.analysis.setup_ac_paper_observation --weekly >> /app/logs/setup_ac_paper_obs_$(date +\%Y\%m\%d)_weekly.log 2>&1

# --- DROPPED (PR-B): host crons that point at missing/stale scripts — NOT migrated ---
#   phase2_daily_verification  — scripts/cron/phase2_daily_verification.sh is MISSING (dead cron)
#   counterfactual_weekly      — scripts/cron/counterfactual_weekly.sh is MISSING (dead cron)
#   forecast_weekly_report     — forecast_vs_rl_comparison.py is RL-removed/stale
# PR-D host-cron teardown removes these dead entries from the host crontab.
